# podcastindex-proxy/deploy/launchd/de.mediathekview.podcastindex-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.mediathekview.podcastindex-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/podcastindex-proxy</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PODCASTINDEX_API_KEY</key>
        <string>hier-key-eintragen</string>
        <key>PODCASTINDEX_API_SECRET</key>
        <string>hier-secret-eintragen</string>
        <key>PODCASTINDEX_PROXY_HOST</key>
        <string>127.0.0.1</string>
        <key>PODCASTINDEX_PROXY_PORT</key>
        <string>18080</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/usr/local/var/podcastindex-proxy</string>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/podcastindex-proxy.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/podcastindex-proxy.log</string>
</dict>
</plist>
